# deploy/crime-cartography-inbox-worker.service.example
[Unit]
Description=Crime Cartography read-only mailbox intake
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=crime-cartography
Group=crime-cartography
WorkingDirectory=/srv/crime-cartography
Environment=CRIME_CARTOGRAPHY_GMAIL_CLIENT_SECRET=/etc/crime-cartography/gmail-oauth-client.json
Environment=CRIME_CARTOGRAPHY_GMAIL_TOKEN=/var/lib/crime-cartography/gmail-readonly-token.json
ExecStart=/usr/bin/node pipeline/subscribers/list-subscriptions.mjs
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/crime-cartography

[Install]
WantedBy=multi-user.target
